# deploy/telegram-bot.service
[Unit]
Description=Claude Superpowers Telegram Bot
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=ray
WorkingDirectory=/home/ray/claude-superpowers
EnvironmentFile=/home/ray/claude-superpowers/.env
ExecStart=/home/ray/claude-superpowers/.venv/bin/python -m telegram_bot.entrypoint
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
